# pkg/osafi/osafi-desktop.install
post_install() {
    post_upgrade
}

post_upgrade() {
    ln --force --symbolic /usr/bin/google-chrome-stable /usr/bin/google-chrome

    systemctl enable --now lightdm
}

pre_remove() {
    rm --force /usr/bin/google-chrome
}
