# pkg/osafi/osafi-base.install
post_install() {
    post_upgrade
}

post_upgrade() {
    ln -sf /usr/share/zoneinfo/America/Los_Angeles /etc/localtime
    echo "LANG=en_US.UTF-8" > /etc/locale.conf
    sed 's/#en_US/en_US/' -i /etc/locale.gen
    locale-gen

    systemctl enable --now nordvpnd
    systemctl enable --now docker
}
