# pkg/osafi/osafi-dell-xps.install
post_install() {
    post_upgrade
}

post_upgrade() {
    systemctl enable --now lightdm.server
}
